# bakup-auto.service
[Unit]
Description="run /usr/local/bin/bakup"

[Service]
Environment="STARTED_BY_SYSTEMD=yes"
ExecStart=/usr/local/bin/bakup --run
